# roles/httpd/templates/httpd_supervisord.conf.j2
[program:httpd]
directory={{ httpd_server_root }}
command={{ httpd_server_root }}/bin/httpd -C "ErrorLog /dev/stdout" -DFOREGROUND
stdout_logfile={{ httpd_server_root }}/logs/httpd.log
stderr_logfile={{ httpd_server_root }}/logs/httpd_err.log
stdout_events_enabled=true
stderr_events_enabled=true
autorestart=true
stopasgroup=true
killasgroup=true
